# data/unity8.conf
description "Unity Shell v8"
author "Ricardo Mendoza <ricmm@ubuntu.com>"

start on ((xsession SESSION=ubuntu-touch) or (xsession SESSION=ubuntu-touch-surfaceflinger)) and started dbus
stop on desktop-end

respawn

script
    # Point all future jobs in this session to our Mir socket instead of
    # unity-system-compositor's socket.
    export MIR_SERVER_FILE=$XDG_RUNTIME_DIR/mir_socket
    initctl set-env --global MIR_SOCKET=$MIR_SERVER_FILE
    exec ${BINARY:-unity8} $ARGS
end script
