# data/unity8-dash.conf
description "Unity Shell v8 Dash"
author "Michael Zanetti <michael.zanetti@canonical.com>"

emits scope-ui-starting

start on started unity8
stop on stopping unity8

kill timeout 120

# Focused app is 100 so we want to be less likely to get
# killed than the focused app but more than Unity8

oom score 50

respawn
respawn limit unlimited

env APP_ID=unity8-dash

pre-start script
    initctl emit scope-ui-starting
end script

exec ${BINARY:-unity8-dash} $ARGS --desktop_file_hint=/usr/share/applications/unity8-dash.desktop
