# data/unity8-greeter-ofono.conf
description "Unity 8 Greeter ofono start"
author "Michael Terry <mterry@ubuntu.com>"

start on unity8-greeter-started
task

script
    start --no-wait ofono-setup
end script
