# data/unity8-dash.conf
description "Unity Shell v8 Dash"
author "Michael Zanetti <michael.zanetti@canonical.com>"

emits scope-ui-starting

start on started unity8
stop on stopping unity8

kill timeout 120

respawn

pre-start script
    if [ -z "$UNITY_SCOPES_LIST" ]; then
        # FIXME: remove once we have this in dconf
        initctl set-env UNITY_SCOPES_LIST="scopes;clickscope;musicaggregator;videoaggregator"
    fi

    initctl emit scope-ui-starting
end script

exec ${BINARY:-unity8-dash} $ARGS --desktop_file_hint=/usr/share/applications/unity8-dash.desktop
