# data/unity8-setcap.conf
# unity8-setcap - ugly hacks to arrange for /usr/bin/unity8 to have CAP_SYS_RESOURCE

# XXX replace me with some root-helper to gain CAP_SYS_RESOURCE XXX

author "Loïc Minier <loic.minier@ubuntu.com>"
description "Ugly hacks to arrange for /usr/bin/unity8 to have CAP_SYS_RESOURCE"

# start when first boot-hooks event is emitted and before lightdm (lightdm
# starts ubuntu-touch-session which starts unity8); note that /run is
# guaranteeds to be be there because lightdm starts on filesystem
start on boot-hooks and starting lightdm

# NOT a task as otherwise this would block restarting lightdm

env RUNDIR=/run/unity8-setcap

# work needs to be done in pre-start as this really is a job with nothing to
# start
pre-start script
    if [ ! -e "$RUNDIR" ]; then
        mkdir "$RUNDIR"
        # /run is noexec, hence mounting another tmpfs exec
        # NB: unity8 is 35K; 512K should be enough for everyone
        mount -o rw,nosuid,nodev,exec,relatime,mode=755,size=512k -t tmpfs tmpfs "$RUNDIR"
        cp -a /usr/bin/unity8 "$RUNDIR"
        setcap CAP_SYS_RESOURCE=+ep "$RUNDIR/unity8"
        # bind-mount this back as unity8 checks dirname(argv[0]'s) == /usr to
        # decide whether it's installed or not
        mount --bind "$RUNDIR/unity8" /usr/bin/unity8
    fi
end script

