# data/unity8.conf
description "Unity Shell v8"
author "Ricardo Mendoza <ricmm@ubuntu.com>"

start on ((xsession SESSION=ubuntu-touch) or (xsession SESSION=ubuntu-touch-surfaceflinger)) and started dbus
stop on desktop-end

respawn

pre-start script
    # Point unity8 at unity-system-compositor
    initctl set-env MIR_SERVER_HOST_SOCKET=$MIR_SOCKET
    # Point all future jobs in this session to our Mir socket instead of
    # unity-system-compositor's socket.
    initctl set-env --global MIR_SOCKET=$XDG_RUNTIME_DIR/mir_socket
end script

exec ${BINARY:-unity8} --file=$MIR_SOCKET $ARGS
