# data/unity8-greeter-ofono.conf
description "Unity 8 Greeter ofono start"
author "Michael Terry <mterry@ubuntu.com>"

start on unity8-greeter-started
stop on started :sys:ofono
task

script
    # Is ofonod ready?
    if pidof onofod >/dev/null; then
        stop
    else
        # Wait 30 seconds.  If ofono starts earlier, upstart will kill us.
        # If 30 seconds go by, try starting ofono-setup anyway, but it will
        # probably fail.
        sleep 30
    fi
end script

post-stop script
    start --no-wait ofono-setup
end script
