# data/unity8-greeter-init.conf
description "Unity 8 Greeter initialization"
author "Michael Terry <mterry@ubuntu.com>"

start on unity8-greeter-started
task
emits indicator-services-start

script
    start --no-wait ofono-setup || true
    initctl emit --no-wait indicator-services-start || true
    start --no-wait maliit-server || true
end script
