# data/unity8.conf
description "Unity Shell v8"
author "Ricardo Mendoza <ricmm@ubuntu.com>"

start on ((xsession SESSION=ubuntu-touch) or (xsession SESSION=ubuntu-touch-surfaceflinger)) and started dbus
stop on desktop-end

respawn

exec unity8 ${BINARY:-unity8} $ARGS
