# data/unity8-greeter-init.conf
description "Unity 8 Greeter initialization"
author "Michael Terry <mterry@ubuntu.com>"

start on unity8-greeter-started
emits indicator-services-start
task

script
    initctl start --no-wait ofono-setup || true
    initctl start --no-wait maliit-server || true
end script
