# data/unity8-greeter-logrotate.conf
description "Unity 8 Greeter logrotate support"
author "Michael Terry <mterry@ubuntu.com>"

start on unity8-greeter-started
task

script
    # logrotate normally starts on startup and then sleeps for 60 seconds.
    # Otherwise, it skips that sleep.  Since we can't emit the startup signal
    # for a unity8-greeter session, we have to trigger logrotate manually.
    # So we do the sleep ourselves in its stead.
    sleep 60
    initctl emit :sys:rotate-logs
end script
