# data/unity8-filewatcher.conf
description "File system watcher for unity8"
author "Michael Zanetti <michael.zanetti@canonical.com>"

# This method of updating the launcher will not work in our new snappy world,
# neither for snaps installed in the system nor if we are running inside a snap.
# Plus it hardcodes the phablet user...  To be replaced once bug 1630997 is
# fixed in ubuntu-app-launch.

# Workaround for bug 1360208. ~ and * in one expression doesn't work currently
start on (file FILE=/home/phablet/.local/share/applications/*.desktop) or (file FILE=/usr/share/applications/*.desktop)

exec dbus-send --type=method_call --dest=com.canonical.Unity.Launcher /com/canonical/Unity/Launcher com.canonical.Unity.Launcher.Refresh

